# docker/etc/service/nginx/run
#!/bin/sh
set -e

sv start pm2 || exit 1

exec /usr/sbin/nginx -c /etc/nginx/nginx.conf;