# docker/etc/services.d/pm2/run
#!/usr/bin/with-contenv sh

cd /home/app
HOME=/home/app exec su -c "pm2-docker start process.yml" app
